# package/debian/control/prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/jddclient" ]; then
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d jddclient stop || true
	else
		/etc/init.d/jddclient stop || true
	fi
fi